# luci-app-waytix/root/etc/crontabs/root
# Обновление списка серверов каждый день в 3:00
0 3 * * * /etc/waytix/update.sh >/dev/null 2>&1

# Проверка подключения каждые 5 минут
*/5 * * * * /etc/init.d/waytix check >/dev/null 2>&1
